# systemd/.config/systemd/user/update-wallpaper.service
[Unit]
Description=Download and set a random wallpaper

[Service]
Type=oneshot
ExecStart=%h/.local/bin/wallpaper update --quiet
