# systemd/.config/systemd/user/update-wallpaper.timer
[Unit]
Description=Update wallpaper periodically

[Timer]
OnBootSec=2min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
